# changesmbpasswd.socket
[Unit]
Description=Change smb passwd scgi Socket

[Socket]
ListenStream=4000

[Install]
WantedBy=sockets.target
