# changesmbpasswd.service
[Service]
ExecStart=/usr/local/bin/changesmbpasswd 4000

[Install]
WantedBy=multi-user.target
